# agent-runtime/orchestration/systemd/agent-runtime.service
# agent-runtime.service — one agent activation cycle via Docker.
#
# Managed by agent-runtime.timer (fires every 5 minutes).
# Type=oneshot so systemd waits for the container to finish before
# marking the unit active — this pairs correctly with the timer's
# OnUnitActiveSec setting.
#
# Installation:
#   sudo cp agent-runtime.service agent-runtime.timer /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now agent-runtime.timer
#
# Logs:
#   journalctl -u agent-runtime.service -f

[Unit]
Description=Agent Runtime — single activation cycle
Documentation=https://github.com/tiendv89/agent-workflow/tree/main/agent-runtime
After=docker.service
Requires=docker.service

[Service]
Type=oneshot

# Load secrets from a protected env file.
# Create: sudo install -m 600 -o root /dev/null /etc/agent-runtime/env
# Then add the following variables:
#   ANTHROPIC_API_KEY=sk-ant-...
#   GITHUB_TOKEN=ghp_...        # GitHub PAT with repo scope — required by pr-create skill
#   GIT_AUTHOR_EMAIL=agent@example.com
#   GIT_AUTHOR_NAME=Agent Bot
#
# SSH key — systemd EnvironmentFile does NOT support multi-line values,
# so a PEM key cannot be pasted directly as SSH_PRIVATE_KEY.
# Recommended: use SSH_KEY_PATH pointing to a key file on disk.
#   SSH_KEY_PATH=/etc/agent-runtime/ssh/id_rsa
# The key file must be mode 400, readable by the user running the service.
#
# Alternative (advanced): encode the key on one line and decode in the container.
#   SSH_PRIVATE_KEY=$(awk 'NF {sub(/\r/, ""); printf "%s\\n",$0;}' ~/.ssh/id_rsa)
# Then the container runtime detects the escape sequences and reformats the key.
# This requires runtime support — not yet implemented; use SSH_KEY_PATH instead.
EnvironmentFile=/etc/agent-runtime/env

# Agent configuration lives at /etc/agent-runtime/agent.yaml
# SSH key lives at /etc/agent-runtime/ssh/id_rsa  (mode 400)

ExecStart=/usr/bin/docker run --rm \
    --name agent-runtime \
    -e ANTHROPIC_API_KEY \
    -e GITHUB_TOKEN \
    -e GIT_AUTHOR_NAME \
    -e GIT_AUTHOR_EMAIL \
    -e SSH_KEY_PATH=/agent/ssh/id_rsa \
    -v /etc/agent-runtime/agent.yaml:/agent/agent.yaml:ro \
    -v /etc/agent-runtime/ssh:/agent/ssh:ro \
    -v agent-runtime-workspaces:/agent/data \
    ghcr.io/tiendv89/agent-runtime:latest

# Pull the latest image before each run so the host stays current.
ExecStartPre=/usr/bin/docker pull ghcr.io/tiendv89/agent-runtime:latest

StandardOutput=journal
StandardError=journal
SyslogIdentifier=agent-runtime

# The agent exits 0 on idle cycles and on kill-switch. Non-zero exits
# (exit 2/3/4) are recorded in the journal but do not trigger a restart
# (that is the timer's job).
SuccessExitStatus=0 2 3 4

[Install]
WantedBy=multi-user.target
